# scripts/com.statesmedia.wiki-lint.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.statesmedia.wiki-lint</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>wiki-lint-cron.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>scripts-path-placeholder</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>log-path-placeholder/wiki-lint-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>log-path-placeholder/wiki-lint-stderr.log</string>
</dict>
</plist>
